# systemd/.config/systemd/user/mutt-sync.service
[Unit]
Description=Sync mail and update notmuch index

[Service]
Type=oneshot
ExecStart=%h/.mutt/sync-and-index.sh
