# systemd/.config/systemd/user/mutt-sync.timer
[Unit]
Description=Run mail sync every 10 minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=10m
Unit=mutt-sync.service

[Install]
WantedBy=timers.target
